# .config/systemd/user/noctalia.service
[Unit]
Description=Noctalia shell
After=graphical-session.target
PartOf=graphical-session.target
StartLimitIntervalSec=0

[Service]
Type=simple
Environment=XDG_CURRENT_DESKTOP=Hyprland
Environment=XDG_SESSION_TYPE=wayland
ExecStartPre=/bin/sh -lc 'test -n "${WAYLAND_DISPLAY:-}" && test -S "$XDG_RUNTIME_DIR/$WAYLAND_DISPLAY"'
ExecStart=/usr/bin/qs -c noctalia-shell
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
